# supervisor.conf
[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
user=root
autostart=true
stdout_events_enabled=true
stderr_events_enabled=true